# systemd/kazoo.target
[Unit]
Description=Started kazoo daemons containers
After=docker.service

[Install]
WantedBy=multi-user.target
